# .openode.cron
# m h  dom mon dow   command
5 */6 * * * cd /opt/app/ && forever restartall