# .openode.cron
# m h  dom mon dow   command
5 */6 * * * cd /opt/app/ && forever restartall
5 1 * * * find /opt/data/logs/*.log -mtime +7 -delete